# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=systemctl start crond.service
